# build/crontab
0 0 * * * /usr/local/bin/renew_certs.sh >> /var/log/cron.log 2>&1
